# recipes-core/1u0022-setup-hostname/1u0022-setup-hostname/sethostname.service
[Unit]
Description=Set Hostname
After=basic.target
Wants=basic.target

[Service]
Type=oneshot
ExecStart=/usr/sbin/set_hostname.sh

[Install]
WantedBy=multi-user.target
