# recipes-core/1u0022-setup-hostname/1u0022-setup-hostname/sethostname.service
[Unit]
Description=Set Hostname
DefaultDependencies=no
RequiresMountsFor=/srv
Before=sysinit.target

[Service]
Type=oneshot
ExecStart=/usr/sbin/set_hostname

[Install]
WantedBy=multi-user.target
